# rdup.cron
# sample cron file of rdup
11 02 * * * root rdup-simple /etc /home /var file:///raid/backup/$HOSTNAME 
