# rdup.cron
########################################
# cron file of rdup                    #
########################################
 
# --- dump - each new month a full dump--- #
10 02 * * * root /usr/sbin/dump.sh -b /vol/backup/elektron/`date +%Y%m` -b /vol/raid/backup/elektron/`date +%Y%m` /etc /home /tftpboot
00 02 * * * root /usr/sbin/dump.sh -b /vol/backup/elektron/`date %Y%m` -b /vol/raid/backup/elektron/`date +%Y%m` /var/mail /var/svn /var/lib
